# deploy/auto-fund-poll.timer
[Unit]
Description=Poll fund blogger RSS every 30 minutes (legacy; prefer auto-fund-daily.timer)

[Timer]
OnBootSec=5min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
