# deploy/auto-fund-poll.timer
[Unit]
Description=Poll fund blogger RSS at 14:30 and 14:35 on weekdays

[Timer]
OnCalendar=Mon-Fri *-*-* 14:30:00
OnCalendar=Mon-Fri *-*-* 14:35:00
Persistent=true

[Install]
WantedBy=timers.target
