# deploy/auto-fund-poll.timer
[Unit]
Description=Poll fund blogger RSS at 14:20 and 14:25 on weekdays

[Timer]
OnCalendar=Mon-Fri *-*-* 14:20:00
OnCalendar=Mon-Fri *-*-* 14:25:00
Persistent=true

[Install]
WantedBy=timers.target
