# deploy/auto-fund-digest.timer
[Unit]
Description=Send fund digest at 14:40 on weekdays (trading day check inside job)

[Timer]
OnCalendar=Mon-Fri *-*-* 14:40:00
Persistent=true

[Install]
WantedBy=timers.target
