# deploy/auto-fund-daily.timer
[Unit]
Description=Weekday 14:20 poll RSS then send digest (trading day check in digest)

[Timer]
OnCalendar=Mon-Fri *-*-* 14:20:00
Persistent=true

[Install]
WantedBy=timers.target
